# shsh-deploy/roles/shsh/files/upstart/shsh-api-pictures.conf
start on runlevel [2345]
stop on runlevel [!2345]

respawn
respawn limit 10 5

setuid root
chdir /opt/shsh-api

env NODE_ENV="production"
env SHSH_SECRET="shshcl4v3"

exec npm run pictures